# templates/etc/init/docker-run.conf.erb
# This file is managed by Puppet and local changes
# may be overwritten

description "start and stop <%= @title %> in docker"
author "Gareth Rushgrove"

start on (started <%= @start_on %>)
stop on (stopping <%= @stop_on %>)

respawn
respawn limit 5 20

pre-start script
  <% if @verify_checksum -%>
  # require file saying checksum was verified OK
  # (file created by /usr/local/bin/verify_docker_image, invoked by Puppet docker::image)
  test -f /var/lib/puppet/state/verify_docker_image.<%= @sanitised_image %> || { stop; exit 0; }
  <% end -%>

  <% if @pre_start -%>
     env ACTION="pre-start" \
         CID="`cat /var/run/docker-<%= @sanitised_title %>.cid 2>/dev/null`" \
         CIDFILE="/var/run/docker-<%= @sanitised_title %>.cid" \
	 IMAGE="<%= @image %>" \
         NAME="<%= @name %>" \
         <%= @pre_start %>
  <% end -%>
  true
end script

post-start script
   <% if @post_start -%>
     env ACTION="post-start" \
         CID="`cat /var/run/docker-<%= @sanitised_title %>.cid 2>/dev/null`" \
         CIDFILE="/var/run/docker-<%= @sanitised_title %>.cid" \
         IMAGE="<%= @image %>" \
	 NAME="<%= @name %>" \
         <%= @post_start %>
   <% end -%>
   <% @trigger_restart_array.each do |image| -%>
   service docker-<%= image %> restart
   <% end -%>
   true
end script

script
<% if @use_name -%>
  <%= @docker_command %> inspect <%= @name %> && <%= @docker_command %> start -a -i <%= @name %>
  <%= @docker_command %> inspect <%= @name %> || exec <%= @docker_command %> run --cidfile=/var/run/docker-<%= @sanitised_title %>.cid \
<% else -%>
  exec <%= @docker_command %> run --cidfile=/var/run/docker-<%= @sanitised_title %>.cid \
<% end -%>
    <% if @username %>-u '<%= @username %>' <% end %> \
    <% if @hostname %>-h '<%= @hostname %>'<% end %> \
    <% if @disable_network %>-n false<% end %> \
    <% if @privileged %>--privileged<% end %> \
    <% if @dns %><% @dns_array.each do |address| %> \
        --dns <%= address %><% end %> \
    <% end %> \
    <% if @env %><% @env_array.each do |env| %> \
        -e <%= env %><% end %> \
    <% end %> \
    <% if @ports %><% @ports_array.each do |port| %> \
        -p <%= port %><% end %> \
    <% end %> \
    <% if @expose %><% @expose_array.each do |expose_port| %> \
        --expose=<%= expose_port %><% end %> \
    <% end %> \
    <% if @volumes %><% @volumes_array.each do |volume| %> \
        -v <%= volume %><% end %> \
    <% end %> \
    <% if @volumes_from %> --volumes-from <%= @volumes_from %><% end %> \
    <% if @net %>--net <%= @net %><% end %> \
    <% if @lxc_conf_array %><% @lxc_conf_array.each do |lxcconf| %> \
        --lxc-conf="<%=lxcconf%>"<% end %> \
    <% end %> \
    -m <%= @memory_limit %> \
    <% if @links %><%= @links_array.map {|link| \
        "--link #{link}" }.join(' ') %> \
    <% end %> \
    <% if @use_name %>--name <%= @name %><% end %> \
    <% if @extra_parameters %><% @extra_parameters_array.each do |param| %>  \
      <%= param %> <% end %> \
    <% end -%> \
    <%= @image %> \
    <% if @command %> <%= @command %><% end %>
end script

post-stop script
  <% if @post_stop -%>
   env ACTION="post-stop" \
       CID="`cat /var/run/docker-<%= @sanitised_title %>.cid 2>/dev/null`" \
       CIDFILE="/var/run/docker-<%= @sanitised_title %>.cid" \
       IMAGE="<%= @image %>" \
       NAME="<%= @name %>" \
       <%= @post_stop %>
  <% end -%>
  if [ -e "/var/run/docker-<%= @sanitised_title %>.cid" ]; then
      <%= @docker_command %> kill "$(cat /var/run/docker-<%= @sanitised_title %>.cid)" && \
      rm "/var/run/docker-<%= @sanitised_title %>.cid"
  fi
  true
end script

pre-stop script
  <% if @pre_stop -%>
     env ACTION="pre-stop" \
         CID="`cat /var/run/docker-<%= @sanitised_title %>.cid 2>/dev/null`" \
         CIDFILE="/var/run/docker-<%= @sanitised_title %>.cid" \
         IMAGE="<%= @image %>" \
	 NAME="<%= @name %>" \
	 <%= @pre_stop %>
  <% end -%>
  true
end script
